# conf.d/facerec.conf
[program:facerec]
command=/usr/bin/startservices.sh start
user=root
autostart=true
autorestart=true
startretries=3
